# templates/etc/cron.d/atop
* * * * *   root /usr/local/bin/fix_atop.sh
